# data/dock.service.in
[Unit]
Description=Pantheon Dock
PartOf=graphical-session.target
After=gnome-session-initialized.target

[Service]
Slice=session.slice
Type=dbus
BusName=io.elementary.dock
ExecStart=@bindir@/io.elementary.dock

# Lower down dock's OOM score to avoid being killed by OOM-killer too early
OOMScoreAdjust=-1000

[Install]
WantedBy=gnome-session@pantheon.target
